# meta-adv-imx/recipes-utils/gester/gester/gester.service
[Unit]
Description=a gesture interpreting program for touchscreens 

[Service]
ExecStart=/usr/bin/gester

[Install]
WantedBy=multi-user.target
